# update-home-ip.timer
[Unit]
Description=Update Home IP

[Timer]
OnCalendar=*-*-* *:37:42

[Install]
WantedBy=timers.target
